# services/.config/systemd/user/backup.service
[Unit]
Description=Backup to pi in dorm

[Service]
Type=simple
ExecStart=/home/alex/bin/backup

[Install]
WantedBy=multi-user.target
