# services/.config/systemd/user/selfsniper.service
[Unit]
Description=Self Sniper

[Service]
Type=simple
ExecStart=/home/alex/bin/selfsniper
StandardError=journal
